# deploy/systemd/fun-agent-service.service
[Unit]
Description=fun-agent-service (Node + Python)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/fun-ai-station/fun-agent-service
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=2

# 不依赖 systemd 环境变量：server.js 会读取：
#   fun-agent-service/config/fun-agent-service.env

[Install]
WantedBy=multi-user.target

